# plotting/plot-starter.service
[Unit]
Description=Start Plotter
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/plot_starter
#ExecStart=/usr/local/bin/plot_starter --config /etc/chiagarden/default.config

Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
